# scripts/librarian.service
# scripts/librarian.service
# Deploy: cp to ~/.config/systemd/user/
[Unit]
Description=Assay Librarian Agent
After=network.target

[Service]
Type=oneshot
EnvironmentFile=/home/morgan/librarian/.env
WorkingDirectory=/home/morgan/librarian
ExecStart=/home/morgan/librarian/venv/bin/python /home/morgan/assay/scripts/librarian.py
TimeoutStartSec=600
